# myHomeSim.service
[Unit]
Description=myHomeSim node.js
Documentation=https://example.com
After=network.target mysgw.service

[Service]
Environment=NODE_PORT=8080 NODE_ENV=production
Type=simple
User=pi
Group=pi
ExecStart=/usr/bin/node /home/pi/myHomeSim/index.js

Restart=on-failure
# Restart service after 10 seconds if node service crashes
RestartSec=10

# Output to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=myHomeSimServer

[Install]
WantedBy=multi-user.target
